# amazon-s3-backup/rootfs/etc/services.d/amazon-s3-backup/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# Home Assistant Community Add-on: Amazon S3 Backup
# ==============================================================================
s6-svscanctl -t /var/run/s6/services