# amazon-s3-backup/rootfs/etc/services.d/amazon-s3-backup/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Amazon S3 Backup
# ==============================================================================

bashio::log.info "Starting Amazon S3 Backup..."

bucket_name="$(bashio::config 'bucket_name')"
storage_class="$(bashio::config 'storage_class')"
bucket_region="$(bashio::config 'bucket_region')"
monitor_path="/backup"

export AWS_ACCESS_KEY_ID="$(bashio::config 'aws_access_key')"
export AWS_SECRET_ACCESS_KEY="$(bashio::config 'aws_secret_access_key')"
export AWS_REGION="$bucket_region"

bashio::api.supervisor GET /supervisor/backups false '.data.backups|=sort_by(.date)'

bashio::log.info "Using AWS CLI version: '$(aws --version)'"
bashio::log.info "Command: 'aws s3 sync $monitor_path s3://$bucket_name/ --no-progress --region $bucket_region --storage-class $storage_class'"
aws s3 sync $monitor_path s3://$bucket_name/ --no-progress --region $bucket_region --storage-class $storage_class

bashio::log.info "Finished Amazon S3 Backup."